# files/reposync.service
[Unit]
Description=Update local repository mirrors

[Service]
Type=simple
ExecStart=/etc/scripts/reposync-all.sh
